# atelier-backend/atelier-backend.service
[Unit]
Description=Atelier Backend API
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/atelier-backend
ExecStart=/usr/bin/node src/index.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=atelier-backend

# Environment
EnvironmentFile=/var/www/atelier-backend/.env

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
